# claude-proxy.service
[Unit]
Description=Claude API Proxy
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/claude-proxy
EnvironmentFile=/opt/claude-proxy/.env
ExecStart=/opt/claude-proxy/venv/bin/python /opt/claude-proxy/start.py
Restart=on-failure
RestartSec=5
StandardInput=tty
TTYPath=/dev/tty1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
